# OnDevice/mcStart.service
[Unit]
Description=MC_Server_BootUP

[Service]
ExecStart=/home/florianwimbauer/raspMC/OnDevice/boot.sh
Restart=on-failure
RemainAfterExit=yes
Type=forking

[Install]
WantedBy=multi-user.target