# OnDevice/mcStart.service
[Unit]
Description=Minecraft Server

[Service]
Type=simple
ExecStart=/usr/bin/screen -DmS minecraft /home/florianwimbauer/raspMC/OnDevice/boot.sh
WorkingDirectory=/opt/minecraft-server
Restart=on-failure

[Install]
WantedBy=multi-user.target