# OnDevice/mcStart.service
[Unit]
Description=MC_Server_BootUP

[Service]
ExecStart=/home/florianwimbauer/boot.sh
Restart=no

[Install]
WantedBy=multi-user.target